# com.jewhatredtoday.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jewhatredtoday.update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/filitti/Downloads/jew-hatred-today/run.sh</string>
    </array>

    <!-- Run every 4 hours (14400 seconds), and once right after loading. -->
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/filitti/Downloads/jew-hatred-today/run.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/filitti/Downloads/jew-hatred-today/run.log</string>
</dict>
</plist>
